# ft_services/srcs/influxDB/configs/supervisord.conf
[supervisord]
nodaemon=true
[program:telegraf]
command= sh -c "telegraf && kill -s SIGTERM $(cat supervisord.pid)"
[program:influxdb]
command= sh -c "influxd run -config /etc/influxdb.conf && kill -s SIGTERM $(cat supervisord.pid)"